# NZBmegasearch/large_server/supervisor/conf.d/detachedtrend.conf
[program:detached]
command = /usr/bin/python detachedtrend.py
directory = /opt/usntssearch/NZBmegasearch/
user = zerobyte
autostart = true
autorestart = true
