# NZBmegasearch/large_server/supervisor/conf.d/mega.conf
[program:mega]
command = /usr/local/bin/gunicorn mega2:app --bind unix:/tmp/gunicorn.sock -w 15
directory = /opt/usntssearch/NZBmegasearch/
user = zerobyte
autostart = true
autorestart = true
