# NZBmegasearch/large_server/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
user = root
autostart = true
autorestart = true
